# arch-update.service
[Unit]
Description=Execute arch-updates

[Service]
Type=oneshot

ExecStart=/usr/bin/arch_update
